# hotpants.service
[Unit]
Description=NYT R&D HOT_PANTS PROX SENSOR

[Service]
Type=simple
WorkingDirectory=/home/root/hotpants/
ExecStart=/usr/bin/python hotpants-VCNL4000.py
KillMode=process
