# hotpants.timer
[Unit]
Description=NYT R&D HOT_PANTS timer

[Timer]
OnBootSec=200
OnUnitInactiveSec=10
Unit=hotpants.service